# vps-container/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/tmp/supervisord.pid
user=root

[program:xvfb]
command=/usr/bin/Xvfb :1 -screen 0 %(ENV_VNC_GEOMETRY)sx%(ENV_VNC_DEPTH)s -ac +extension GLX +render -noreset
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/xvfb.log
stderr_logfile=/var/log/supervisor/xvfb.err

[program:xfce]
command=/usr/bin/dbus-run-session -- /usr/bin/startxfce4
environment=DISPLAY=":1",HOME="/home/dyad",USER="dyad"
user=dyad
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/supervisor/xfce.log
stderr_logfile=/var/log/supervisor/xfce.err

[program:x11vnc]
command=/usr/local/bin/start-x11vnc.sh
autostart=true
autorestart=true
priority=30
stdout_logfile=/var/log/supervisor/x11vnc.log
stderr_logfile=/var/log/supervisor/x11vnc.err

[program:novnc]
command=/usr/bin/websockify --web=/usr/share/novnc/ 127.0.0.1:6080 127.0.0.1:5900
autostart=true
autorestart=true
priority=40
stdout_logfile=/var/log/supervisor/novnc.log
stderr_logfile=/var/log/supervisor/novnc.err

[program:filebrowser]
command=/usr/local/bin/start-filebrowser.sh
autostart=true
autorestart=true
priority=45
stdout_logfile=/var/log/supervisor/filebrowser.log
stderr_logfile=/var/log/supervisor/filebrowser.err

[program:dyad]
command=/usr/local/bin/start-dyad.sh
autostart=true
autorestart=true
priority=60
stdout_logfile=/var/log/supervisor/dyad.log
stderr_logfile=/var/log/supervisor/dyad.err

[program:nginx]
command=/usr/local/bin/start-nginx.sh
autostart=true
autorestart=true
priority=90
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx.err
